# ot-recorder.service
[Unit]
Description=OT Video Recorder
After=network-online.target
Wants=network-online.target
# Do not Requires=dev-video0.device — if the camera is missing at boot,
# systemd times out that device unit and stays dead even after hot-plug.
# The recorder waits for the camera in-process instead.

[Service]
Type=simple
User=pi
Group=pi

EnvironmentFile=/etc/ot-agent/env

ExecStart=/usr/local/bin/ot-recorder

Restart=always
RestartSec=5

# Give ffmpeg time to flush its current segment on SIGTERM
TimeoutStopSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ot-recorder

# Prevent runaway disk usage
LimitFSIZE=10G

[Install]
WantedBy=multi-user.target
